# scripts/firefox-commander.service
[Unit]
Description=Launch Firefox in kiosk mode for Stream Commander
# Ensure stream-overlord is up first
After=stream-overlord.service
Wants=stream-overlord.service

# Start after your desktop session is ready
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Delay before starting (change 5 to whatever)
ExecStartPre=/usr/bin/sleep 5

ExecStart=/usr/bin/firefox --kiosk --no-remote --new-instance \
  http://localhost:8105/commander

Restart=on-failure
RestartSec=5s
Environment=NODE_ENV=production

[Install]
WantedBy=graphical-session.target
